# .config/systemd/user/cruft.target
[Unit]
Description=Stuff after login
#Wants=multi-user.target

[Install]
Alias=default.target
